# root/etc/cron.d/permission
15 */4 * * * root /.pyenv/shims/python /scripts/permission.py
