# root/etc/cron.d/permission
15 */4 * * * root /usr/local/bin/python /scripts/permission.py
