# root/etc/cron.d/permission
15 */2 * * * root /.pyenv/shims/python /scripts/permission.py
